# supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/opt/mailgw.sh

[program:rsyslog]
command=/usr/sbin/rsyslogd -n -c3
